# deploy/gex-dashboard.service
[Unit]
Description=GEX Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/meshulro/Projects/gex-dashboard
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/meshulro/Projects/gex-dashboard/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
